# config/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid
user=root

[program:slmodemd]
priority=10
command=/bin/bash -lc "ulimit -n 1024; exec slmodemd -e /usr/local/bin/slmodem-sip-bridge"
autostart=true
autorestart=true
startsecs=2
startretries=3
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:oob-hub]
priority=20
command=/usr/local/bin/oob-hub
autostart=true
autorestart=true
startsecs=2
startretries=3
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
